# backend/srvus.service
[Unit]
After=network.target

[Service]
ExecStart=/usr/local/bin/srvus
ReadOnlyDirectories=/
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
